# docker-next/s6/.s6-svscan/finish
#!/bin/sh

# Cleanup s6 event folder
rm -rf "$(find /app/gogs/docker/s6/ -name 'event')" 2>/dev/null || true
