# docker/s6/gogs/run
#!/bin/sh

if test -f ./setup; then
    source ./setup
fi

export USER=git
exec gosu $USER /app/gogs/gogs web
