# docker/etc/profile.d/gitea.sh
#!/bin/bash
export GITEA_CUSTOM=/data/gitea
