# docker/s6/syslogd/run
#!/bin/sh

if test -f ./setup; then
    source ./setup
fi

if [ "$(grep -Ei 'debian' /etc/*release)" ]; then
    exec gosu root /sbin/syslogd -nS
else
    exec gosu root /sbin/syslogd -nS -O-
fi

