# docker/s6/syslogd/run
#!/bin/sh

if test -f ./setup; then
    source ./setup
fi

exec gosu root /sbin/syslogd -nS -O-
