# docker/s6/syslogd/run
#!/bin/sh

if test -f ./setup; then
    # shellcheck disable=SC2039,SC1091,SC3046
    source ./setup
fi

# Syslogd typically requires root. If not running as root, skip or use alternative logging.
if [ "$(id -u)" != "0" ]; then
    echo "syslogd:disabled | Running as non-root user, syslogd disabled." 1>&2
    exec sleep infinity
fi

exec /sbin/syslogd -nS -O-
