# docker/s6/otelcol/run
#!/bin/sh

exec /app/otelcol/otelcol-contrib --config /app/gogs/docker/templates/config.yaml