# docker/etc/s6/openssh/run
#!/bin/bash
[[ -f ./setup ]] && source ./setup

pushd /root > /dev/null
    exec su-exec root /usr/sbin/sshd -D
popd
