# docker/s6/openssh/run
#!/bin/sh

# SSH service is disabled by default in rootless mode.
# To enable SSH, the container must run as root or with appropriate capabilities.
# Consider using Gogs' built-in SSH server instead by configuring START_SSH_SERVER = true in app.ini

# Check if we can bind to port 22 (requires root or CAP_NET_BIND_SERVICE)
if [ "$(id -u)" != "0" ]; then
    echo "ssh:disabled| Running as non-root user, SSH daemon disabled. Use built-in SSH server instead." 1>&2
    exec sleep infinity
fi

if test -f ./setup; then
    # shellcheck disable=SC2039,SC1091,SC3046
    source ./setup
fi

exec /usr/sbin/sshd -D -f /app/gogs/docker/sshd_config
