# scripts/init/gentoo/gogs
#!/sbin/openrc-run

DIR=/home/git/gogs
USER=git
PORT=3000

start_stop_daemon_args="--user ${USER} --chdir ${DIR}"
command="${DIR}/gogs"
command_args="web -port ${PORT}"
command_background=yes
pidfile=/var/run/gogs.pid

depend()
{
    need net
}
